# scripts/etc/systemd/system/brightness_mon.service
[Unit]
Description=Backlight monitor

[Service]
ExecStart="/usr/local/bin/backlightmon"

[Install]
WantedBy=multi-user.target
